# dev-kit/build-system/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=<mail>
HOME=/home/<user>

10 21 * * * /home/shared/GNOME/gnome-live/dev-kit/build-system/cron.sh > /home/<user>/cronlog
